# ngrok.service
[Unit]
Description=Ngrok

[Service]
ExecStart=/home/said/bin/ngrok http 3000

[Install]
WantedBy=multi-user.target

